# supervisor-laravel-queue.conf.example
# Laravel Queue Worker Supervisor 配置示例
# 
# 使用方法：
# 1. 复制此文件到 /etc/supervisor/conf.d/laravel-queue.conf
# 2. 修改路径和用户
# 3. 运行: sudo supervisorctl reread
# 4. 运行: sudo supervisorctl update
# 5. 运行: sudo supervisorctl start laravel-queue-resource-push:*
# 6. 运行: sudo supervisorctl start laravel-queue-ota-push:*

[program:laravel-queue-resource-push]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/your/project/artisan queue:work redis --queue=resource-push --tries=3 --timeout=300 --sleep=3 --max-jobs=1000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/path/to/your/project/storage/logs/resource-push.log
stopwaitsecs=3600

[program:laravel-queue-ota-push]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/your/project/artisan queue:work redis --queue=ota-push --tries=3 --timeout=300 --sleep=3 --max-jobs=1000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/path/to/your/project/storage/logs/ota-push.log
stopwaitsecs=3600

